# config/.config/systemd/user/chromium.service
[Unit]
Description=Chromium

[Service]
ExecStart=/usr/bin/chromium
Environment=DISPLAY=:0
Restart=on-failure

[Install]
WantedBy=default.target
